# reflector/reflector.service
[Unit]
Description=Pacman mirrorlist update

[Service]
Type=oneshot
ExecStart=/usr/bin/reflector --protocol http --latest 50 --sort rate --fastest 10 --save /etc/pacman.d/mirrorlist.reflector
