# systemd/autostart-xdg.service
[Unit]
Description=Autostart XDG

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/dex --autostart
