# systemd/chromium.service
[Unit]
Description=Chromium
After=xorg.target

[Service]
ExecStart=/usr/bin/chromium

[Install]
WantedBy=stuff.target
