# systemd/chromium.service
[Unit]
Description=Chromium

[Service]
ExecStart=/usr/bin/chromium

[Install]
WantedBy=default.target
