# systemd/dropbox.service
[Unit]
Description=Dropbox

[Service]
ExecStart=/usr/bin/sh "%h/.dropbox-dist/dropboxd"

[Install]
WantedBy=default.target
