# systemd/dropbox.service
[Unit]
Description=Dropbox

[Service]
ExecStart=/usr/bin/env "%h/.dropbox-dist/dropboxd"
