# systemd/firefox@.service
[Unit]
Description=Firefox
After=network.target

[Service]
ExecStart=/usr/bin/firefox -P %I -no-remote

[Install]
WantedBy=main-stuff.target
