# systemd/firefox.service
[Unit]
Description=Firefox
After=network.target

[Service]
ExecStart=/usr/bin/firefox

[Install]
WantedBy=main-stuff.target
