# systemd/firefox.service
[Unit]
Description=Firefox

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/firefox
