# systemd/gaming-stuff.target
[Unit]
Description=Target for xinitrc stuff
Wants=gaming.target

[Install]
Alias=gaming.target
