# systemd/git-gc@.service
[Unit]
Description=GC git repo %f

[Service]
ExecStart=/usr/bin/git -C %f gc --aggressive
