# systemd/git-gc@.timer
[Unit]
Description=Run git gc periodically

[Timer]
OnCalendar=Tue,Wed,Thu,Fri
