# systemd/gvim.service
[Unit]
Description=GVim
After=xorg.target

[Service]
ExecStart=/usr/bin/vim -gf

[Install]
WantedBy=stuff.target
