# systemd/gvim.service
[Unit]
Description=GVim

[Service]
ExecStart=/usr/bin/vim -gf

[Install]
WantedBy=default.target
