# systemd/gvim.service
[Unit]
Description=GVim

[Service]
ExecStart=/usr/bin/vim -gf
Restart=always

[Install]
WantedBy=main-stuff.target
