# systemd/lxpolkit.service
[Unit]
Description=LXDE policykit authenticator

[Service]
ExecStart=/usr/bin/lxpolxit
