# systemd/lxpolkit.service
[Unit]
Description=LXDE policykit authenticator

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/lxpolkit
