# systemd/main-stuff.target
[Unit]
Description=Target for startup programs
Wants=main.target
