# systemd/mpd.service
[Unit]
Description=Music Player Daemon
After=sound.target network.target

[Service]
ExecStart=/usr/bin/mpd --no-daemon
Restart=always
RestartSec=2
