# systemd/nasa-sdo-wallpaper.timer
[Unit]
Description=Update the NASA SDO wallpaper

[Timer]
OnBootSec=0
OnUnitInactiveSec=30min
