# systemd/nasa-sdo-wallpaper.timer
[Unit]
Description=Hourly update the NASA SDO wallpaper

[Timer]
OnBootSec=0
OnUnitInactiveSec=1h
