# systemd/pidgin.service
[Unit]
Description=Pidgin
After=network.target

[Service]
Environment=DISPLAY=:0 BROWSER=/usr/bin/firefox
ExecStart=/usr/bin/pidgin
