# systemd/skype.service
[Unit]
Description = Skype

[Service]
ExecStart = /usr/bin/skype --disable-api

[Install]
WantedBy = default.target
