# systemd/skype.service
[Unit]
Description=Skype
After=xorg.target network.target
Wants=pulseaudio.service

[Service]
EnvironmentFile=/home/daniel/.systemd-environment
ExecStart=/usr/bin/skype --disable-api

[Install]
WantedBy=main-stuff.target
