# systemd/steam.service
[Unit]
Description=Steam
After=network.target

[Service]
ExecStart=/usr/bin/steam

[Install]
WantedBy=gaming-stuff.target
