# systemd/steam.service
[Unit]
Description=Steam
After=network.target

[Service]
Environment=DISPLAY=:1 BROWSER=/usr/bin/firefox
ExecStart=/usr/bin/steam
