# systemd/unclutter.service
[Unit]
Description=Unclutter

[Service]
ExecStart=/usr/bin/unclutter
Restart=always

[Install]
WantedBy=main-stuff.target
