# systemd/xmonad.service
[Unit]
Description=Xmonad
Before=main-stuff.target

[Service]
ExecStart=/usr/bin/xmonad
Restart=always
RestartSec=10

[Install]
WantedBy=main.target
