# systemd/.config/systemd/user/git-gc@.service
[Unit]
Description=Garbage collect git repository in %I

[Service]
ExecStart=/usr/bin/git -C %I gc
