# systemd/.config/systemd/user/redshift.service
[Unit]
Description=Redshift

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/redshift -m randr -c %h/.config/redshift
