# systemd/.config/systemd/user/stand-up.timer
[Unit]
Description=Stand up occasionally

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
