# systemd/.config/systemd/user/stand-up.service
[Unit]
Description=Stand up notification

[Service]
ExecStart=/usr/bin/notify-send "Stand up!" --icon=up
