# systemd/.config/systemd/user/steam.service
[Unit]
Description=Steam

[Service]
Environment=DISPLAY=:1
ExecStart=/usr/bin/steam
