# systemd/.config/systemd/user/twmnd.service
[Unit]
Description=twmnd

[Service]
ExecStart=/usr/bin/twmnd
