# systemd/.config/systemd/user/wallpaper.timer
[Unit]
Description=Change the wallpaper periodically

[Timer]
OnBootSec=0
OnUnitInactiveSec=30min
