# systemd/.config/systemd/user/weechat.service
[Unit]
Description=WeeChat

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/urxvt -c irc -e "/usr/bin/weechat-curses"
