# wallpaper/random-wallpaper.timer
[Unit]
Description=Run random-wallpaper daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
